# debian/heat-api.init.in
#!/bin/sh
### BEGIN INIT INFO
# Provides:          heat-api
# Required-Start:    $network $local_fs $remote_fs $syslog
# Required-Stop:     $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      0 1 6
# Short-Description: Heat API
# Description:       Heat API server
### END INIT INFO

# Author: Thomas Goirand <zigo@debian.or>

# PATH should only include /usr/* if it runs after the mountnfs.sh script
DESC="Heat API"
PROJECT_NAME=heat
NAME=${PROJECT_NAME}-api
